# install/systemd/mssh-server.service
[Unit]
Description=mssh rendezvous server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mssh server --host 0.0.0.0 --port 8443
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
